# rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user start-up script.
#

ntpd -s

exit 0
